# i3lock-on-lid-close.service
/etc/systemd/system/

[Unit]
Description=i3lock on Lid Close
Before=sleep.target

[Service]
User=your_username
Type=forking
Environment=DISPLAY=:0
ExecStart=/usr/bin/i3lock

[Install]
WantedBy=sleep.target

sudo systemctl enable i3lock-on-lid-close.service
sudo systemctl start i3lock-on-lid-close.service
sudo systemctl status i3lock-on-lid-close.service
